# duckdb-ui.service
[Unit]
Description=DuckDB UI - basedosdados explorer
After=network.target

[Service]
Type=simple
User=YOUR_USER
WorkingDirectory=/path/to/baseldosdados
ExecStartPre=/usr/bin/python3 prepara_gui.py
ExecStart=/usr/bin/duckdb --ui basedosdados.duckdb
Restart=on-failure
RestartSec=5s
EnvironmentFile=/path/to/baseldosdados/.env

[Install]
WantedBy=multi-user.target
